# systemd/airband-ui.service
[Unit]
Description=Airband Gain/Squelch Web UI
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
EnvironmentFile=-/etc/airband-ui.conf
ExecStart=/usr/bin/python3 /opt/airband-ui/airband_ui.py
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
