# systemd/rtl-airband.service
[Unit]
Description=RTLSDR-Airband (scan -> Icecast)
After=network-online.target icecast2.service
Wants=network-online.target
Requires=icecast2.service
Wants=rtl-airband-last-hit.service

[Service]
Type=simple
Environment=CONFIG_DIR=/usr/local/etc
Environment=CONFIG_SYMLINK=/usr/local/etc/rtl_airband.conf
Environment=GROUND_CONFIG_PATH=/usr/local/etc/rtl_airband_ground.conf
Environment=COMBINED_CONFIG_PATH=/usr/local/etc/rtl_airband_combined.conf
EnvironmentFile=-/etc/airband-ui.conf
ExecStartPre=-/usr/bin/pkill -x rtl_airband
ExecStartPre=/usr/bin/python3 /home/willminkoff/scannerproject/scripts/build-combined-config.py
ExecStart=/bin/bash /home/willminkoff/scannerproject/scripts/rtl-airband-filter.sh ${COMBINED_CONFIG_PATH}
Restart=on-failure
RestartSec=0
TimeoutStopSec=1
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
