# systemd/rtl-airband.service
[Unit]
Description=RTLSDR-Airband (scan -> Icecast)
After=network-online.target icecast2.service
Wants=network-online.target
Requires=icecast2.service
Wants=rtl-airband-last-hit.service

[Service]
Type=simple
ExecStartPre=-/usr/bin/pkill -x rtl_airband
ExecStartPre=/usr/bin/python3 /home/willminkoff/scannerproject/scripts/build-combined-config.py
ExecStart=/usr/local/bin/rtl_airband -F -e -c /usr/local/etc/rtl_airband_combined.conf
Restart=on-failure
RestartSec=0
TimeoutStopSec=1
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
