# systemd/rtl-airband.service
[Unit]
Description=RTLSDR-Airband (scan -> Icecast)
After=network-online.target icecast2.service
Wants=network-online.target
Requires=icecast2.service

[Service]
Type=simple
ExecStartPre=-/usr/bin/pkill -x rtl_airband
ExecStart=/home/willminkoff/scannerproject/scripts/rtl-airband-with-freq.sh
Restart=on-failure
RestartSec=2
TimeoutStopSec=3

[Install]
WantedBy=multi-user.target
