# systemd/rtl-airband-last-hit.service
[Unit]
Description=RTLSDR-Airband last-hit tracker
After=rtl-airband.service rtl-airband-ground.service
Wants=rtl-airband.service rtl-airband-ground.service
PartOf=rtl-airband.service rtl-airband-ground.service

[Service]
Type=simple
ExecStart=/home/willminkoff/scannerproject/scripts/rtl-airband-last-hit.sh
Restart=always
RestartSec=1

[Install]
WantedBy=rtl-airband.service rtl-airband-ground.service
