# systemd/rtl-airband-last-hit.service
[Unit]
Description=RTLSDR-Airband last-hit tracker
After=rtl-airband.service
Requires=rtl-airband.service
PartOf=rtl-airband.service

[Service]
Type=simple
ExecStart=/bin/bash /home/willminkoff/scannerproject/scripts/rtl-airband-last-hit.sh
Restart=always
RestartSec=1

[Install]
WantedBy=rtl-airband.service
