# systemd/sb3-ap-fallback.service
[Unit]
Description=SB3 AP Fallback (start AP when LAN unreachable)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/sb3-ap-fallback.conf
ExecStart=/usr/local/sbin/sb3-ap-fallback
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
