# systemd/scanner-digital-mixer.service
[Unit]
Description=Digital Audio Mixer (RTL-Airband + SDRTrunk -> Icecast)
After=network-online.target icecast2.service rtl-airband.service scanner-digital.service
Wants=network-online.target icecast2.service rtl-airband.service scanner-digital.service
Requires=icecast2.service
StartLimitIntervalSec=0

[Service]
Type=simple
EnvironmentFile=-/etc/airband-ui.conf
ExecStart=/usr/bin/python3 /home/willminkoff/scannerproject/scripts/digital-mixer.py
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
